# app/conf.d/fluxbox.conf
[program:fluxbox]
command=fluxbox
autorestart=true
user=root