# app/conf.d/xterm.conf
[program:xterm]
command=xterm
autorestart=true
user=burp